# packaging/systemd/arctis-sound-manager.service
# arctis-sound-manager.service — systemd user unit for the ASM daemon
#
# Install (user-level, no root required):
#   systemctl --user enable --now arctis-sound-manager.service
#
# After the AppImage or .deb/.rpm is installed the daemon binary is at
# /usr/bin/asm-cli (system-wide install) or
# ~/.local/bin/asm-cli (per-user AppImage install).
# Adjust ExecStart below to match the installed location.

[Unit]
Description=Arctis Sound Manager daemon
Documentation=https://github.com/oxibux/arctis-sound-manager
After=pipewire.service wireplumber.service
Wants=pipewire.service wireplumber.service

[Service]
Type=simple
# OWNER-FILL: adjust the path if the binary is installed elsewhere.
# For AppImage-extracted installs, this is typically:
#   ExecStart=%h/.local/bin/asm-cli daemon
# For system-wide .deb/.rpm installs:
#   ExecStart=/usr/bin/asm-cli daemon
ExecStart=/usr/bin/asm-cli daemon
Restart=on-failure
RestartSec=3s
# Inherit the user's DBUS/PipeWire environment so asm-cli can reach PipeWire.
Environment=XDG_RUNTIME_DIR=%t
# Runtime dir for the Unix socket (default: $XDG_RUNTIME_DIR/arctis-sound-manager.sock)
RuntimeDirectory=arctis-sound-manager
RuntimeDirectoryMode=0700

[Install]
WantedBy=default.target
